# test/rootfs/overlay/etc/systemd/system/foo.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/foo

[Install]
WantedBy=multi-user.target
